# configs/systemd/user/x11off.target
[Unit]
Description=X11 is off
Requires=default.target
Conflicts=x11on.target
After=default.target
